# configs/automatic-backup/kiwitcms-deployment-backup.service
[Unit]
Description=Kiwi TCMS Deployment Backup

[Service]
ExecStart=/usr/bin/kiwitcms-deployment-backup.sh
